# xcel-itron-mqtt/rootfs/etc/s6-overlay/s6-rc.d/xcel-itron-mqtt/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Xcel iTron2MQTT
# Runs Xcel iTron2MQTT
# ==============================================================================
bashio::log.info "Starting Xcel iTron2MQTT"

# https://github.com/zaknye/xcel_itron2mqtt/tree/main#options
export MQTT_SERVER=$(bashio::config mqtt "host")
export MQTT_USER=$(bashio::config mqtt "username")
export MQTT_PASSWORD=$(bashio::config mqtt "password")
export MQTT_PORT=$(bashio::config mqtt "port")
export METER_IP=$(bashio::config 'meter_ip')
export METER_PORT=$(bashio::config 'meter_port')

# Persistent Storage
export CERT_DIR=$(bashio::config 'cert_dir')
export CERT_PATH=${CERT_DIR}/$(bashio::config 'cert_file')
export KEY_PATH=${CERT_DIR}/$(bashio::config 'key_file')

# Use openssl.conf to enable unsafe cert usage
export OPENSSL_CONF=/opt/xcel_itron2mqtt/openssl.conf

cd /opt/xcel_itron2mqtt || return
exec python3 -Wignore main.py
